# linux/root/etc/s6-overlay/s6-rc.d/runner/finish
#!/bin/bash

echo "$1" > /run/s6-linux-init-container-results/exitcode
exec /run/s6/basedir/bin/halt
