# linux/root/etc/s6-overlay/s6-rc.d/runner/run
#!/command/execlineb -P
s6-setuidgid runner
/var/runner/runner.sh
