# linux/root/etc/s6-overlay/s6-rc.d/docker/run
#!/command/execlineb -P
s6-notifyoncheck -c "docker ps"
dind dockerd
